# gist-manager-updater.service
[Unit]
Description=Update Gist Manager from origin/main
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/gist-manager
Environment=HOME=/root
Environment=BUN_INSTALL_CACHE_DIR=/var/cache/gist-manager-updater/bun
Environment=XDG_CACHE_HOME=/var/cache/gist-manager-updater
EnvironmentFile=/etc/gist-manager/gist-manager-updater.env
ExecStart=/opt/gist-manager/scripts/update-production.sh
TimeoutStartSec=20min
CacheDirectory=gist-manager-updater
CacheDirectoryMode=0700

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
ReadWritePaths=/opt/gist-manager
